# Original_Update/update_HLS-T_1.17.tgz/update-files.tar.gz/etc/hotplug.d/default/photomodel.hotplug
#!/bin/sh -x
if /usr/local/bin/isphotomodel.sh; then
	/usr/local/bin/hotplug.sh $*
fi
